# debian/postinst
#!/bin/sh
# postinst script for #PACKAGE#
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        APPCONF=/etc/ampweb/production.ini
        APACHECONF=/etc/ampweb/ampweb.conf
        CACERT=/usr/lib/python2.7/dist-packages/ampweb/static/cacert.pem

        if [ ! -f $APPCONF ]; then
            cp /etc/ampweb/production.ini-dist $APPCONF
            SECRET=`cat /dev/urandom | tr -dc 'a-zA-Z0-9' | fold -w 16 | head -n 1`
            sed -i "s/auth.secret = CHANGEME/auth.secret = $SECRET/" $APPCONF

            # if netevmon was already installed then enable events now
            if dpkg -s netevmon 2>/dev/null | \
                    grep -qx "Status: .* installed"; then
                sed -i \
                    "s/ampweb.disableevents = yes/ampweb.disableevents = no/" \
                    /etc/ampweb/production.ini
            fi
        fi

        # create the ampweb user if it doesn't already exist, so we can run
        # commands as them rather than www-data
        if ! getent passwd ampweb > /dev/null; then
            adduser --quiet --system --group --no-create-home \
                    --home /etc/ampweb/ ampweb
        fi

        # if amppki is already installed then link in the cacert
        if dpkg -s amppki 2>/dev/null | \
                grep -qx "Status: .* installed"; then
            if [ ! -e $CACERT ]; then
                ln -s /etc/amppki/cacert.pem $CACERT
            fi
        fi

        # TODO use dh_apache2
        if [ -d /etc/apache2/conf-available/ ]; then
            # try to install the config for apache2 in jessie or similar
            if [ ! -e /etc/apache2/conf-available/ampweb.conf ]; then
                ln -s $APACHECONF /etc/apache2/conf-available/
                /usr/sbin/a2enconf -q ampweb
            fi
        elif [ -d /etc/apache2/conf.d/ ]; then
            # try to install the config for apache2 in wheezy or similar
            if [ ! -e /etc/apache2/conf.d/ampweb.conf ]; then
                ln -s $APACHECONF /etc/apache2/conf.d/
            fi
        else
            echo "ERROR: Couldn't determine apache2 configuration directory"
            exit 1
        fi

        # always restart, as either config changed or website files changed
        if apache2ctl configtest 2>/dev/null; then
            #invoke-rc.d apache2 force-reload || true
            # use "restart" as a workaround for bug #558608
            invoke-rc.d apache2 restart || true
        else
            echo "Your apache2 configuration is broken, please fix it and restart apache2 manually."
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
